# supervisord.conf
[supervisord]
nodaemon = true

[program:cron]
command=sh /entry.sh
user=root
autostart=true
autorestart=true

[program:mosdns]
command=sh /entrypoint.sh
user=root
autostart=true
autorestart=true
